# lol-bot.service
[Unit]
Description=League of Legends Discord Stats Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/lol-discord-bot
ExecStart=/home/YOUR_USERNAME/lol-discord-bot/venv/bin/python bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment (alternatively, use EnvironmentFile)
# EnvironmentFile=/home/YOUR_USERNAME/lol-discord-bot/.env

[Install]
WantedBy=multi-user.target
